# root/etc/init.d/nora
#!/bin/sh /etc/rc.common

. /lib/functions.sh

USE_PROCD=1
START=95
STOP=10

SERVICE_NAME="nora"

start_service() {
	local install_dir

	config_load nora
	config_get install_dir config install_dir '/opt/nora'
	if [ ! -x "$install_dir/bin/nora" ]; then
		logger -t nora "Nora binary is not installed: $install_dir/bin/nora"
		return 1
	fi
	/usr/libexec/nora-control render-nora-config >/dev/null || return 1

	procd_open_instance
	procd_set_param env NORA_CONFIG_PATH="$install_dir/config/config.toml"
	procd_set_param command /usr/libexec/nora-run
	procd_set_param respawn 3600 5 5
	procd_close_instance
}

stop_service() {
	uci -q delete firewall.nora >/dev/null 2>&1 || true
	uci commit firewall >/dev/null 2>&1 || true

	if [ -x /etc/init.d/firewall ]; then
		/etc/init.d/firewall reload >/dev/null 2>&1 || /etc/init.d/firewall restart >/dev/null 2>&1 || true
	fi
}

reload_service() {
	stop
	start
}

service_triggers() {
	procd_add_reload_trigger "nora"
}
